# anyone-stick-cm-warmcheck.service
[Unit]
Description=Anyone Stick CM Warmcheck

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cm_warmcheck.sh
